# examples/systemd/rest-server.socket
[Socket]
ListenStream = 8000

[Install]
WantedBy = sockets.target
