# examples/bsd/openbsd
#!/bin/ksh
#
# $OpenBSD: $

daemon="/usr/local/bin/rest-server"
daemon_flags="--path /var/restic"
daemon_user="_restic"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
